# touchpad/disable-touchpad.service
# move me to /etc/systemd/system

[Unit]
Description=Disable touchpad

[Service]
ExecStart=modprobe -r psmouse

[Install]
WantedBy=multi-user.target
